# docker/custom/supervisor/conf.d/temporal-ui.conf
[program:temporal-ui]
command=/container/bin/start-temporal-ui.sh
directory=/etc/temporal
user=root
autostart=true
autorestart=true
priority=25
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
